# Scripts/rc.local
/usr/sbin/wifi_connect.sh
/usr/sbin/wlan0_ip_add.sh
/usr/sbin/mount_ssd.sh
exit 0